# sources/bt_hfp_chatd_c_android4.4_v3.3.4/service/linux/bt_chatd.service
[Unit]
Description=bt_chatd (Bluetooth HFP chat daemon)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/bt_chatd --adapter-index 0 --tcp 127.0.0.1:3333 --audio linuxraw --pcm-capture /dev/snd/pcmC0D0c --pcm-playback /dev/snd/pcmC0D0p
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
